# alpine/alpine-ssh.service.in
[Unit]
Description=Alpine Developer Environment with local SSH Access
Documentation=man:sshd(8) man:sshd_config(5)

[Service]
Type=forking
Restart=always
ExecStartPre=/usr/bin/rm -f %t/%n-pid %t/%n-cid
ExecStart=%%PODMAN%% run \
    --replace \
    --conmon-pidfile %t/%n-pid \
    --cidfile %t/%n-cid \
    --publish %%PORT%%:%%PORT%% \
    --volume %h:/home/alpine \
    --security-opt label=disable \
    --userns=keep-id \
    --tty \
    --name alpine-development-environment \
    --detach \
    alpine-dev-env
ExecStop=%%PODMAN%% stop \
    --time 2 \
    --cidfile %t/%n-cid
ExecStopPost=%%PODMAN%% rm \
    --ignore \
    --force \
    --cidfile %t/%n-cid
KillMode=none
PIDFile=%t/%n-pid

[Install]
WantedBy=multi-user.target
